# root/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bash

cd /config

exec \
	s6-setuidgid abc mosquitto -c /config/mosquitto.conf
